# files/60-onbattery.rules
SUBSYSTEM=="power_supply",RUN+="/usr/local/bin/ac_changed.sh"
